# src/homie-nodes/homie_camnode.service
[Unit]
Description=Camnode Test Server

[Service]
ExecStart=mktemp
Restart=always
RestartSec=60

[Install]
WantedBy=default.target
